# data/hp-omen-key.service
[Unit]
Description=HP Omen Key Handler — launches HP Manager on Omen Key press
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/libexec/hp-manager/omen-key-listener.sh
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
